# Aposa2024-raspberrypi-main/systemd/auto2.service
[Unit]
Description=mabiki

[Service]
Type=simple
ExecStart=/bin/bash /root/auto2.sh

[Install]
WantedBy=multi-user.target
